# usr/lib/systemd/user/cachyos-gamescope-autologin.service
[Unit]
Description=Set autologin to gamescope on logout
PartOf=graphical-session.target
After=graphical-session.target
ConditionPathExists=!%t/gamescope-environment

[Service]
Type=oneshot
ExecStart=steamos-session-select gamescope --no-logout

[Install]
WantedBy=graphical-session.target
